# log_collector.service
[Unit]
Description=EZREC Log Collector Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=michomanoly14892
Group=michomanoly14892
WorkingDirectory=/opt/ezrec-backend
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/ezrec-backend
ExecStart=/usr/bin/python3 /opt/ezrec-backend/log_collector.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=log_collector

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/ezrec-backend

[Install]
WantedBy=multi-user.target 